# etc/systemd/system/nextcloud-cron.timer
[Unit]
Description=Run NextCloud background Jobs.

[Timer]
OnCalendar=*:0,15,30,45
Persistent=true

[Install]
WantedBy=timers.target
